# puppet-dashboard/configs/supervisor/apache2.conf
[program:apache2]
priority=10
directory=/tmp
command=/usr/sbin/apache2ctl -D FOREGROUND
user=root
autostart=true
autorestart=true
