# puppet-dashboard/configs/supervisor/puppet-dashboard.conf
[program:puppet-dashboard]
priority=30
directory=/app
command=bundle exec rails server
environment=RAILS_ENV="production"
user=root
autostart=true
autorestart=true
